# .config/systemd/user/desk.service
[Unit]
Description=Manage my Idasen desk
[Service]
ExecStart=/home/fra/.local/bin/cron/desk
