# .config/systemd/user/hide-popup.service
[Unit]
Description=Hide screen sharing popups
[Service]
ExecStart=/home/fra/.local/bin/cron/hide-popup
